# white_list_4/regexdns.service
[Unit]
Description=Regex-based DNS filter proxy
After=network.target dnsmasq.service
Wants=dnsmasq.service

[Service]
ExecStart=/opt/regexdns/venv/bin/python /opt/regexdns/regexdns.py
Environment=PYTHONPATH=/opt/regexdns/venv/lib/python3.13/site-packages
Restart=always
RestartSec=5
User=root
Group=root
WorkingDirectory=/opt/regexdns

[Install]
WantedBy=multi-user.target
