# auto-login/autologin.service.example
[Unit]
Description=Enable autologin before nightly job

[Service]
Type=oneshot
ExecStart=/usr/bin/sudo /home/<USER>/toggle-autologin.sh