# patrol/patrol.service.example
[Unit]

[Service]
Type=oneshot
ExecStart=/home/${USER}/Linux_device_manager/patrol/patrol.sh
WorkingDirectory=/home/${USER}/Linux_device_manager/patrol
Environment="HOME=/home/${USER}"