# patrol/patrol.service
[Unit]
Description=Run patrol script as root

[Service]
Type=oneshot
ExecStart=/opt/patrol/patrol.sh
WorkingDirectory=/opt/patrol